# lpmd-openrc
#!/sbin/openrc-run

name="lpmd daemon"
description="small linux power management daemon"
command="/usr/bin/lpmd"
command_args=""
pidfile="/run/lpmd.pid"
error_logger="/usr/bin/logger --tag lpmd"
command_background=true
